# test/etc/systemd/system/default.target
[Unit]
Description=My Custom Target
Requires=custom-init.service
After=custom-init.service
AllowIsolate=yes